# sapi/fpm/php-fpm.service.in
[Unit]
Description=The PHP FastCGI Process Manager
After=syslog.target network.target

[Service]
PIDFile=@localstatedir@/run/php-fpm.pid
ExecStart=@sbindir@/php-fpm --nodaemonize --fpm-config @sysconfdir@/php-fpm.conf
ExecReload=/bin/kill -USR2 $MAINPID

[Install]
WantedBy=multi-user.target

